# deploy/systemd/senaran.service.example
[Unit]
Description=Gunicorn Senaran Queue
After=network.target

[Service]
User=root
Group=www-data
WorkingDirectory=/root/WMS/deploy/apps/senaran.cvbjas.com
EnvironmentFile=-/root/WMS/.env.senaran
Environment="APP_TIMEZONE=Asia/Jakarta"
Environment="SENARAN_DATABASE=/var/lib/senaran/antrian.db"
Environment="SESSION_COOKIE_SECURE=true"
RuntimeDirectory=senaran
RuntimeDirectoryMode=0755
StateDirectory=senaran
StateDirectoryMode=0750
ExecStartPre=/usr/bin/python3 /root/WMS/deploy/apps/senaran.cvbjas.com/init_db.py
ExecStart=/usr/bin/python3 /usr/local/bin/gunicorn -w 1 --timeout 180 --graceful-timeout 30 --bind unix:/run/senaran/gunicorn.sock --umask 007 app:app
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
