# deploy/systemd/wms-db-backup.timer.example
[Unit]
Description=Run WMS database backup twice daily (02:00, 12:00)

[Timer]
OnCalendar=*-*-* 02:00:00
OnCalendar=*-*-* 12:00:00
Persistent=true
Unit=wms-db-backup.service

[Install]
WantedBy=timers.target
