# deploy/systemd/wms-db-backup.service.example
[Unit]
Description=WMS database backup job
After=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/root/WMS
EnvironmentFile=-/root/WMS/.env
ExecStart=/usr/bin/python3 /root/WMS/scripts/run_configured_backup.py --output-dir /root/WMS/db_backups --retain-days 14
