# deploy/systemd/wms-db-backup.timer.example
[Unit]
Description=Run WMS SQLite backup 3 times daily (02:00, 14:00, 19:00)

[Timer]
OnCalendar=*-*-* 02:00:00
OnCalendar=*-*-* 14:00:00
OnCalendar=*-*-* 19:00:00
Persistent=true
Unit=wms-db-backup.service

[Install]
WantedBy=timers.target

